# hid-telephony-bridge.service
# systemd user service for the USB HID Telephony mute bridge.
# Runs in the user session alongside PipeWire - no root required
# (assuming the udev rule grants hidraw access).
#
# Install:
#   cp hid-telephony-bridge.service /usr/lib/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now hid-telephony-bridge.service
#
# Logs:
#   journalctl --user -u hid-telephony-bridge -f

[Unit]
Description=USB HID Telephony Mute Bridge
Documentation=https://github.com/fuxx/hid-telephony-bridge
After=pipewire-pulse.service
BindsTo=pipewire-pulse.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/hid-telephony-bridge --verbose
# For a different device, override VID/PID:
# ExecStart=/usr/bin/hid-telephony-bridge --vid 1234 --pid 5678 --verbose
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
