# mv7-hid-bridge.service
# systemd user service for the USB HID Telephony mute bridge.
# Runs in the user session alongside PipeWire - no root required
# (assuming the udev rule grants hidraw access).
#
# Install:
#   cp mv7-hid-bridge.service /usr/lib/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now mv7-hid-bridge.service
#
# Logs:
#   journalctl --user -u mv7-hid-bridge -f

[Unit]
Description=USB HID Telephony Mute Bridge
Documentation=https://github.com/fuxx/mv7-hid-bridge
After=pipewire-pulse.service
BindsTo=pipewire-pulse.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=/usr/bin/mv7-hid-bridge --verbose
# For a different device, override VID/PID:
# ExecStart=/usr/bin/mv7-hid-bridge --vid 1234 --pid 5678 --verbose
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
